# testing/kvm/netbsd/rc.local
#!/bin/sh

if test -r /etc/ifconfig.vioif1 ; then
    echo '==> Hostnames already configured'
elif ifconfig | grep ^vioif1 > /dev/null ; then
    echo '==> Changing hostnames'
    mac=$(ifconfig vioif1 | grep -o -E '([[:xdigit:]]{1,2}:){5}[[:xdigit:]]{1,2}')
    case "$mac" in
	'12:00:00:64:64:23' ) # For East machine
	    hostname=east
	    # eth1
	    vioif1v4=192.1.2.23
	    vioif1v6=2001:db8:1:2::23
	    # eth0
	    vioif2v4=192.0.2.254
	    vioif2v6=2001:db8:0:2::254
	    ;;
	'12:00:00:64:64:45') # For West machine
	    hostname=west
	    # eth1
	    vioif1v4=192.1.2.45
	    vioif1v6=2001:db8:1:2::45
	    # eth0
	    vioif2v4=192.0.1.254
	    vioif2v6=2001:db8:0:1::254
	    ;;
	* )
	    ;;
    esac
    if test -n "${hostname}" ; then
	echo hostname=${hostname}
	for if in vioif1 vioif2 ; do
            v4=$(eval echo \${${if}v4})
            v6=$(eval echo \${${if}v6})
	    echo ${if}v4=${v4}
	    echo ${if}v6=${v6}
	    cat <<EOF > /etc/ifconfig.${if}
inet ${v4} netmask 255.255.255.0
inet6 ${v6} prefixlen 64 alias
EOF
	done
	cat <<EOF >> /etc/ifconfig.vioif1
up
EOF
	echo nameserver 192.1.2.254		>> /etc/resolv.conf
	echo 192.1.2.254			> /etc/mygate
	echo ${hostname}			> /etc/myname
	chmod u=rw,g=r /etc/ifconfig.*
	echo "==> Forcing network configuration"
	hostname ${hostname}
	for if in vioif1 vioif2 ; do
	    while read line ; do
		( set -x ; ifconfig ${if} ${line} )
	    done < /etc/ifconfig.${if}
	done
	echo "==> Done changing hostnames"
    fi
else
    echo "==> Hostnames unchanged"
fi
