# initsystems/systemd/ipsec.service.in
[Unit]
Description=Internet Key Exchange (IKE) Protocol Daemon for IPsec
Documentation=man:libreswan(7) man:ipsec(8) man:pluto(8) man:ipsec.conf(5)
Wants=network-online.target
After=network-online.target

[Service]
Type=@@SD_TYPE@@
NotifyAccess=all
Restart=@@SD_RESTART_TYPE@@
TimeoutStopSec=90s
WatchdogSec=@@SD_WATCHDOGSEC@@
ExecStartPre=@@SBINDIR@@/ipsec checknss
ExecStartPre=@@SBINDIR@@/ipsec checknflog
ExecStart=@@LIBEXECDIR@@/pluto @@SD_PLUTO_OPTIONS@@ --config @@IPSEC_CONF@@ --nofork
ExecStop=@@LIBEXECDIR@@/whack --shutdown
ExecStopPost=@@SBINDIR@@/ipsec stopnflog

[Install]
WantedBy=multi-user.target
