# initsystems/systemd/ipsec.service.in
[Unit]
Description=Internet Key Exchange (IKE) Protocol Daemon for IPsec
Documentation=man:ipsec(8) man:pluto(8) man:ipsec.conf(5)
Wants=network-online.target
After=network-online.target

[Service]
Type=@@SD_TYPE@@
Restart=@@SD_RESTART_TYPE@@

# Set WatchdogSec to the amount of time (in seconds) that systemd will wait
# before restarting an unresponsive pluto.
# EVENT_SD_WATCHDOG updates the heartbeat every 15 seconds, recommended values
# are 60, 90, 120. WatchdogSec=0 disables the action
NotifyAccess=all
WatchdogSec=@@SD_WATCHDOGSEC@@

# Check configuration file
ExecStartPre=@@LIBEXECDIR@@/addconn --config @@IPSEC_CONF@@ --checkconfig
# Check for nss database status and migration
ExecStartPre=@@SBINDIR@@/ipsec checknss
# Check for nflog setup
ExecStartPre=@@SBINDIR@@/ipsec checknflog
# Start the actual IKE daemon
ExecStart=@@LIBEXECDIR@@/pluto @@SD_PLUTO_OPTIONS@@ --config @@IPSEC_CONF@@ --nofork
ExecStop=@@LIBEXECDIR@@/whack --shutdown
ExecStopPost=@@SBINDIR@@/ipsec stopnflog

[Install]
WantedBy=multi-user.target
