# initsystems/systemd/pluto.socket.in
[Unit]
Description=Libreswan IPSec Control Socket

[Socket]
ListenStream=@FINALVARDIR@/run/pluto/pluto.ctl
Service=ipsec.service

[Install]
WantedBy=sockets.target
