# initsystems/upstart/ipsec.conf.in
# ipsec - ipsec services and ike demon (pluto)
#
# IPsec provides encrypted and authenticated communications;
# XFRM is the kernel half of it, Pluto is the user-level
# management daemon.

description	"libreswan ipsec services"

start on runlevel [2345]
stop on runlevel [!2345]

respawn

pre-start script
    cd /
    . /etc/default/pluto
    @@SBINDIR@@/ipsec addconn --config @@IPSEC_CONF@@ --checkconfig
end script

script
    . /etc/default/pluto
    @@SBINDIR@@/ipsec checknss
    @@SBINDIR@@/ipsec checknflog
    @@LIBEXECDIR@@/pluto --config @@IPSEC_CONF@@ --nofork $PLUTO_OPTIONS
    if [ -f @@LIBEXECDIR@@/portexcludes ] ; then
	@@LIBEXECDIR@@/portexcludes
    fi
end script

pre-stop script
    # do not run whack --shutdown in pre-stop, it confused stop into restarting
    # @@SBINDIR@@/ipsec whack --shutdown
    @@SBINDIR@@/ipsec stopnflog
end script
